# corosync.target
[Unit]
Description=Corosync Stack (bare or OpenAIS)

[Install]
WantedBy=multi-user.target
